# deploy/tele-mess-core.user.service
[Unit]
Description=tele-mess-core Telegram archive server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/tele-mess-core
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env tele-mess-core --workspace %h/.local/share/tele-mess-core --config %h/.config/tele-mess-core/config.yml run-server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
