# deploy/tele-mess-core.service
[Unit]
Description=tele-mess-core Telegram archive server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tele-mess-core
WorkingDirectory=/var/lib/tele-mess-core
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/tele-mess-core --workspace /var/lib/tele-mess-core --config /etc/tele-mess-core/config.yml run-server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
